# servicio/servicio.service
[Unit]
Description = Servicio TP

[Service]
Type=simple
Execstart=/bin/bash -c "#reemplazame"

[Install]
WantedBy=multi-user.target


